# daemon/callscreener-daemon.conf
# callscreener-daemon - Call Screener SIP daemon
# Runs as a user session service on Ubuntu Touch
# Place this file at: ~/.config/upstart/callscreener-daemon.conf

description "Call Screener SIP Daemon"
author      "Your Name"

# Start when the user session is ready
start on started dbus

# Restart automatically if it crashes
respawn
respawn limit 5 30

# Run as current user (not root)
setuid $USER

# Working directory
chdir /home/phablet/.local/share/callscreener

# Set up environment
env HOME=/home/phablet
env XDG_RUNTIME_DIR=/run/user/32011
env DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/32011/bus

exec /usr/bin/python3 /home/phablet/.local/share/callscreener/callscreener_daemon.py

# Log output to upstart journal
console log
